# config/systemd/user/swaybg.service
[Unit]
Description=Wallpaper daemon
PartOf=graphical-session.target
After=graphical-session.target
Requisite=graphical-session.target

[Service]
ExecStart=/usr/bin/swaybg -m fill -i "%h/.config/nova/current/background"
Restart=on-failure

[Install]
WantedBy=graphical-session.target
